# content/docs/jfa-go.service
[Unit]
Description=An account management system for Jellyfin.

[Service]
ExecStart=/usr/bin/jfa-go

[Install]
WantedBy=default.target
